# crontab
* * * * * /usr/bin/execute /imapfilter.sh