# crontab
*/5 * * * * /usr/bin/execute /imapfilter.sh